# src/obsidian_sync/templates/com.alexfurrier.obsidian-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alexfurrier.obsidian-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>%HOME%/.local/bin/obsidian-sync</string>
        <string>start</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>%HOME%/.local/state/obsidian-sync/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>%HOME%/.local/state/obsidian-sync/launchd.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>%HOME%/.local/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
